# roles/nodered/templates/backup/backup.timer.j2
[Unit]
Description=Run {{ nodered_service_name }} backup periodically

[Timer]
OnCalendar={{ nodered_backup_cron }}

[Install]
WantedBy=timers.target
